# systemd/borg-drive-handler@.service
[Unit]
Description=Borg backup + livesync handler for external drive (instance: %i)
After=local-fs.target

[Service]
Type=oneshot
# Change to your actual user; do not run this as root.
User=obs1230
Environment=/run/media/obs1230/OliDrive/porta-winux
ExecStart=/home/obs1230/Code/sandbox/porta-winux/scripts/drive_handler.sh
# Prevents overlapping runs if the drive bounces (disconnect/reconnect quickly)
